# centos/redis/redis.service
[Unit]
Description=Redis

[Service]
User=root
ExecStart=/etc/redis/src/redis-server

[Install]
WantedBy=multi-user.target
