# centos/php7-nginx/src/supervisor.d/parse_nginx.conf
[program:parse_nginx]
command=/usr/bin/sh /etc/nginx/conf.d/parse-nginx.sh
autostart=true
autorestart=false
startretries = 1
