# centos/php7-nginx/src/supervisor.d/php_fpm.conf
[program:php_fpm]
command=/usr/sbin/php-fpm -F -R
autostart=true
autorestart=false
numprocs=1
